# roles/dvd_autoplay/files/config/auto-play-dvd.service
[Unit]
Description=Play DVD automatically

[Service]
ExecStart=/usr/local/bin/auto-dvd-mpv.sh
